# deploy/faucet_worker.service
[Unit]
Description=Faucet Worker Automation Service
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/Repositories/cryptobot
EnvironmentFile=/home/azureuser/Repositories/cryptobot/.env
ExecStart=/home/azureuser/Repositories/cryptobot/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=append:/home/azureuser/Repositories/cryptobot/logs/production_run.log
StandardError=append:/home/azureuser/Repositories/cryptobot/logs/production_run.log

[Install]
WantedBy=multi-user.target
